# systemd/cam_snap.timer
[Unit]
Description=Run cam_snap.service every 10 seconds

[Timer]
OnCalendar=*:*:0/10
AccuracySec=1us

[Install]
WantedBy=timers.target
