# systemd/sky_image_aggr-send_storage.service
[Unit]
Description=Uploads images to storage server

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cd /home/pi/Sky-Imager-Aggregator/src && python3 SendStorageV2.py'
User=pi
Group=pi
