# systemd/sky_image_aggr-send_storage.timer
[Unit]
Description=Run sky_image_aggr-send_storage.service hourly

[Timer]
OnCalendar=*:0:0

[Install]
WantedBy=timers.target
